# movie.conf
[program:movie]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/movie
autostart=true
autorestart=true